# docker/config/supervisord_startserver.conf
# Licensed to the Apache Software Foundation (ASF) under one or more
# contributor license agreements.  See the NOTICE file distributed with
# this work for additional information regarding copyright ownership.
# The ASF licenses this file to You under the Apache License, Version 2.0
# (the "License"); you may not use this file except in compliance with
# the License.  You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

[supervisord]
childlogdir=/var/log/supervisord/
logfile=/var/log/supervisord/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
nodaemon=false              ; (start in foreground if true;default false)
pidfile=/var/log/supervisord/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
user=root

[program:nutchserver]
autorestart=true
autostart=true
command=nutch startserver -port %(ENV_SERVER_PORT)s -host %(ENV_SERVER_HOST)s
environment=PATH=/usr/local/bin:%(ENV_PATH)s
process_name=%(program_name)s
numprocs=1
#redirect_stderr=true
startsecs=0
stderr_capture_maxbytes=10MB
stderr_logfile=/var/log/supervisord/%(program_name)s_stderr.log
stderr_logfile_backups=5
stderr_logfile_maxbytes=10MB
#stderr_syslog=
stdout_capture_maxbytes=10MB
stdout_logfile=/var/log/supervisord/%(program_name)s_stdout.log
stdout_logfile_backups=5
stdout_logfile_maxbytes=10MB
#stdout_syslog=
user=root
